# launchd/com.analista.tech.xscan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.analista.tech.xscan</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/pietrosabbatini/Desktop/Analista Tech/.venv/bin/python3</string>
        <string>/Users/pietrosabbatini/Desktop/Analista Tech/xscan.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/pietrosabbatini/Desktop/Analista Tech</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/pietrosabbatini/Desktop/Analista Tech/data/xscan.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/pietrosabbatini/Desktop/Analista Tech/data/xscan.error.log</string>
    <key>RunAtLoad</key><false/>
</dict>
</plist>
