# launchd/com.analista.ngrok.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.analista.ngrok</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ngrok</string>
        <string>http</string><string>8080</string>
        <string>--log</string>
        <string>/Users/pietrosabbatini/Desktop/Analista Tech/data/ngrok.log</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/pietrosabbatini/Desktop/Analista Tech</string>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
</dict>
</plist>
